# scripts/launchd/com.sweepstake.serve.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.sweepstake.serve</string>
  <!-- Serves docs/ + the live-draw API on port 8123. KeepAlive so it comes
       back after a crash or reboot; RunAtLoad so it starts the moment it loads. -->
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>/Users/ianspridgeon/worldcup-sweepstake/scripts/serve.py</string>
    <string>8123</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/ianspridgeon/worldcup-sweepstake</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/ianspridgeon/worldcup-sweepstake/logs/serve.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/ianspridgeon/worldcup-sweepstake/logs/serve.err.log</string>
</dict>
</plist>
