# uk-land-monitoring/deploy/com.landmon.viewer.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- macOS launchd: keep the dashboard running on 127.0.0.1:8000. Replace
     __PROJECT_DIR__. Pair with `tailscale serve 8000` to reach it on your tailnet. -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key> <string>com.landmon.viewer</string>
  <key>ProgramArguments</key>
  <array>
    <string>__PROJECT_DIR__/.venv/bin/python</string>
    <string>-m</string> <string>landmon.web.server</string>
    <string>--host</string> <string>127.0.0.1</string>
    <string>--port</string> <string>8000</string>
  </array>
  <key>WorkingDirectory</key> <string>__PROJECT_DIR__</string>
  <key>EnvironmentVariables</key>
  <dict><key>PYTHONPATH</key> <string>__PROJECT_DIR__/src</string></dict>
  <key>RunAtLoad</key> <true/>
  <key>KeepAlive</key>  <true/>
  <key>StandardOutPath</key>  <string>__PROJECT_DIR__/outputs/viewer.log</string>
  <key>StandardErrorPath</key><string>__PROJECT_DIR__/outputs/viewer.log</string>
</dict>
</plist>
